# ops/docker/transit/crontab
# Transit alerts (watched routes) → Discord, every 5 minutes.
*/5 * * * * python /app/transit_discord.py
